# reverb-ppdb.conf
[program:reverb-ppdb]
command=php /var/www/pesantren-ppdb/artisan reverb:start --host=0.0.0.0 --port=8080
user=www-data
autostart=true
autorestart=true
stderr_logfile=/var/log/reverb-ppdb.err.log
stdout_logfile=/var/log/reverb-ppdb.out.log
